# server/briefing.service
[Unit]
Description=Briefing — personal ops center
After=network.target

[Service]
Type=simple
User=meteopavel
WorkingDirectory=/opt/Briefing
ExecStart=/opt/Briefing/.venv/bin/python -m uvicorn app.web:app --host 127.0.0.1 --port 8001
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
